# hobbyhash-clean/docs/systemd/hobbyhash-wallet.service
[Unit]
Description=HOBC local custodial wallet web app
After=network-online.target hobbyhashd-mainnet.service
Wants=network-online.target hobbyhashd-mainnet.service

[Service]
Type=simple
User=hobbyhashcoin
Group=hobbyhashcoin
WorkingDirectory=/home/hobbyhashcoin/public_html
ExecStart=/usr/bin/php -S 127.0.0.1:18766 -t /home/hobbyhashcoin/public_html
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
